# strava-coach-mcp/ec2-deploy/strava-mcp.service
[Unit]
Description=Strava Coach MCP server
After=network.target

[Service]
WorkingDirectory=/opt/strava-coach-mcp
EnvironmentFile=/opt/strava-coach-mcp/.env
ExecStart=/opt/strava-coach-mcp/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8080
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target
